# ansible/roles/app/templates/dogewand.tmpl.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
exec su -c "/usr/bin/node {{ app.user.home_folder }}/{{ app.name }}/server.js" - {{ app.user.name }}